# files/splunk
#!/bin/sh
#
# /etc/init.d/splunk
# init script for Splunk.
# generated by 'splunk enable boot-start'.
#
# chkconfig: 2345 90 60
# description: Splunk indexer service
#
SPLUNK_HOME="/opt/splunkforwarder"
RETVAL=0

. /etc/init.d/functions

splunk_start() {
  echo Starting Splunk...
  su splunk -c "\"$SPLUNK_HOME/bin/splunk\" start --accept-license --no-prompt --answer-yes"
  RETVAL=$?
}

splunk_stop() {
  echo Stopping Splunk...
  su splunk -c "\"$SPLUNK_HOME/bin/splunk\" stop"
  RETVAL=$?
}

splunk_restart() {
  echo Restarting Splunk...
  su splunk -c "\"$SPLUNK_HOME/bin/splunk\" restart"
  RETVAL=$?
}

splunk_status() {
  echo Splunk status:
  su splunk -c "\"$SPLUNK_HOME/bin/splunk\" status --accept-license --no-prompt --answer-yes"
  RETVAL=$?
}


case "$1" in
  start)
    splunk_start
    ;;
  stop)
    splunk_stop
    ;;
  restart)
    splunk_restart
    ;;
  status)
    splunk_status
    ;;
  *)
    echo "Usage: $0 {start|stop|restart|status}"
    exit 1
esac


exit $RETVAL